# systemd/auto-update.timer
[Unit]
Description=Runs automatic update at 07:00 and 19:00

[Timer]
# Scheduled to run daily at 07:00 and 19:00
OnCalendar=*-*-* 07:00,19:00
Persistent=true

[Install]
WantedBy=timers.target